# btrfs-snapshot.service
[Unit]
Description=Btrfs snapshot

[Service]
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/btrfs-snapshot
