# btrfs-snapshot@.timer
[Unit]
Description=Daily Btrfs snapshot for config %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
