# btrfs-snapshot.service
[Unit]
Description=Btrfs snapshot

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/usr/bin/btrfs-snapshot --confdir
