# systemd/btrfs-snapshot@.timer
[Unit]
Description=Daily Btrfs snapshot on %f

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
